# scripts/backup-timer/kraft-lotto-restore-drill.service
[Unit]
Description=kraft-lotto DB restore drill (quarterly)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=root
EnvironmentFile=/etc/kraft-lotto/backup.env
Environment=DB_HOST=127.0.0.1
ExecStart=/opt/kraft-lotto/scripts/db-restore-drill.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kraft-lotto-restore-drill
